# crontab
1 2 * * * cd /var/www/html/thedailychristorian && ./the_daily_christorian.sh
